# supervisord.conf
[supervisord]
nodaemon=true

[program:teamcity-server]
command=/opt/TeamCity/bin/teamcity-server.sh start

[program:teamcity-agent]
command=/opt/TeamCity/buildAgent/bin/agent.sh run
user=teamcity